# config/com.example.hermes-lab.run-once.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.example.hermes-lab.run-once</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>__REPO_ROOT__/scripts/run-cycle.sh</string>
    <string>--max-runs</string>
    <string>3</string>
  </array>
  <key>StartInterval</key>
  <integer>600</integer>
  <key>EnvironmentVariables</key>
  <dict>
    <key>HERMES_LAB_DATA_ROOT</key>
    <string>__DATA_ROOT__</string>
  </dict>
  <key>StandardOutPath</key>
  <string>/tmp/hermes-lab-run-once.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/hermes-lab-run-once.err</string>
</dict>
</plist>
